# ansible/collections/ansible_collections/ejsmit/containers/roles/nextcloud/templates/nextcloud_cron.service


[Unit]
Description=Run Nextcloud cron job every 5 mins

[Service]
ExecStart=podman exec --user www-data nextcloud php -f /var/www/html/cron.php
